# deepstream-iot/rootfs/etc/services.d/deepstream-iot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: DeepStream IoT Server
# Runs the deepstream app
# ==============================================================================
