# kegbot/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Kegbot Server
# Runs the Redis daemon
# ==============================================================================

bashio::log.info "Starting Redis..."
redis-server
