# kegbot/rootfs/etc/services.d/kegbot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# JPD Hass.io Add-ons: Kegbot Server
# Runs the Kegbot daemon
# ==============================================================================
set -e
CONFIG_PATH=/data/options.json

# Wait for other services to become available
s6-svwait -u -t 5000 /var/run/s6/services/nginx
s6-svwait -u -t 5000 /var/run/s6/services/redis

die() {
  bashio::log.info "Error: $@"
  exit 1
}

do_mysql() {
  mysql -B -u "${KEGBOT_DB_USER}" --password="${KEGBOT_DB_PASS}" -h "${KEGBOT_DB_HOST}" -P ${KEGBOT_DB_PORT} "${@}"
  return $?
}

no_tables() {
  local return_str=$(mysql -B -u "${KEGBOT_DB_USER}" --password="${KEGBOT_DB_PASS}" -h "${KEGBOT_DB_HOST}" -P ${KEGBOT_DB_PORT} "${KEGBOT_DB_NAME}" -e "show tables")
  bashio::log.info "Table_len: ${#return_str}"
  return ${#return_str}
}

wait_for_mysql() {
  nc -z $KEGBOT_DB_HOST $KEGBOT_DB_PORT || sleep 30
  if ! do_mysql "${KEGBOT_DB_NAME}" -e "show tables"; then
    #no db
    do_mysql -e "create database ${KEGBOT_DB_NAME};"
  fi
  if no_tables; then
    #db empty
    /usr/local/bin/pipenv run kegbot migrate --noinput -v 0
    do_mysql "${KEGBOT_DB_NAME}" -e "show tables"
  fi
}

# Perform first-launch setup.
maybe_setup_kegbot() {
  if [ ! -d "${KEGBOT_SETTINGS_DIR}" ]; then
    mkdir -p ${KEGBOT_SETTINGS_DIR}/media
    mv /kegbot-data/local_settings.py ${KEGBOT_SETTINGS_DIR}
  fi
  
  #kegbot collectstatic --noinput -v 0
  #do_mysql -e "create database ${KEGBOT_DB_NAME};" || die "Could not create database."
  true
}

bashio::log.info "Checking environment..."
env
bashio::log.info "Checking Kegbot..."
maybe_setup_kegbot
bashio::log.info "Waiting for Redis..."
redis-cli -h "${KEGBOT_REDIS_HOST}" -p ${KEGBOT_REDIS_PORT} ping
bashio::log.info "Waiting for MySQL..."
wait_for_mysql
bashio::log.info "Starting Kegbot..."
echo `date` >> /kegbot-data/runlog
/usr/local/bin/pipenv run gunicorn pykeg.web.wsgi:application -w 3 -b 127.0.0.1:8000
#kegbot run_all --logs_dir=/kegbot-data --gunicorn_options="-w 3 -b 127.0.0.1:8000"
