# z-way-server/rootfs/etc/services.d/z-way-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# JPD Hass.io Add-ons: Z-Way Server
# Runs the Z-Way Server
# ==============================================================================
# shellcheck disable=SC1091

# Run z-way-server
bashio::log.info 'Starting the Z-Way Server...'
(cd /opt/z-way-server && z-way-server -c /config/z-way-server/config/config.xml -L /dev/stdout)
