# fff-api/rootfs/etc/services.d/fff-api/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# JPD Hass.io Add-ons: FFF-API Server
# Take down the S6 supervision tree when fff-api fails
# ==============================================================================
if -n { s6-test $# -ne 0 }
if -n { s6-test ${1} -eq 256 }

s6-svscanctl -t /var/run/s6/services
