# owntone/rootfs/etc/services.d/owntone/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# JPD Hass.io Add-ons: owntone server
# Runs the owntone Server
# ==============================================================================
# shellcheck disable=SC1091

if ! bashio::fs.directory_exists '/config/owntone/cache'; then
    bashio::log.debug 'Creating cache folder...'
    mkdir -p /config/owntone/cache
fi

if ! bashio::fs.file_exists '/config/forked-daapd/forked-daapd.conf'; then
    bashio::log.debug 'Copying default conf file...'
    cp /usr/local/etc/owntone.conf /config/owntone/owntone.conf
fi

if ! bashio::fs.directory_exists '/config/owntone/music'; then
bashio::log.debug 'Creating music folder...'
    mkdir -p /config/owntone/music
    bashio::log.debug 'Creating HA fifo file...'
    mkfifo -m 666 /config/owntone/music/HomeAssistantAnnounce
fi

bashio::log.info 'Starting the OwnTone Server...'
# Run owntone
owntone -f -c /config/owntone/owntone.conf -P /var/run/owntone.pid
