# forked-daapd/rootfs/etc/services.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# JPD Hass.io Add-ons: forked-daapd server
# Runs the avahi daemon
# ==============================================================================
# shellcheck disable=SC1091

bashio::log.info 'Starting the forked-daapd avahi daemon...'

# Run avahi
exec avahi-daemon --daemonize --no-chroot
