# hass-ap/rootfs/etc/services.d/dhcpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# JPD Hass.io Add-ons: HASS-AP
# Runs dhcpd
# ==============================================================================

true ${INTERFACE:=wlan0}

# start the dhcp server in the foreground
bashio::log.info "Starting DHCP server .."
dhcpd -f ${INTERFACE}
